# profile/profile.d/.bash_logout
# http://lfs.traduc.org/view/blfs-5.1-fr/postlfs/profile.html

# ~/.bash_logout

# C'est un ~/.bash_logout vide qui peut être utilisé comme modèle.
# Vous noterez que la base ~/.bash_logout n'inclut pas de commande clear.
# En effet, celui-ci est géré par le fichier /etc/issue.

# cat > ~/.bash_logout << "EOF"
# Début de ~/.bash_logout
# Écrit pour Beyond Linux From Scratch
# par James Robertson <jameswrobertson@earthlink.net>

# Éléments personnels pour lancer une déconnexion.

# Fin de ~/.bash_logout
#EOF
